# systemd/duckdns.service
[Unit]
Description=DuckDNS Updater

[Service]
Type=oneshot  
ExecStart=/home/vincent/duckdns/duck.sh 

[Install]
WantedBy=multi-user.target  

